# roles/ansible-client/templates/ansible-pull.timer.j2
[Unit]
Description=Starts ansible-pull systemd unit
Requires=ansible-pull.service

[Timer]
Unit=ansible-pull.service
OnCalendar=*-*-* *:00:00
RandomizedDelaySec=600

[Install]
WantedBy=timers.target